# packages/arch/nsweekly.service
[Unit]
Description=NSWeekly service
After=network.target

[Service]
User=nsweekly
Group=nsweekly
ExecStart=/usr/bin/nsweekly
Restart=on-failure

[Install]
WantedBy=multi-user.target
